# state/archlinux/systemd-services/rate-mirrors.timer
[Unit]
Description=Weekly Arch Linux mirrorlist refresh

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
